# services/systemd/knowable-test.service
[Unit]
Description=Knowable Test Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/dmorris/projects/knowable/test
ExecStart=npm run dev
Restart=always
RestartSec=5
Environment=NODE_ENV=development
Environment=PORT=8081
Environment=LOCAL=true
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target