# services/systemd/knowable-frontend.service
[Unit]
Description=Knowable Frontend Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/dmorris/projects/knowable/client
ExecStart=npx serve -p 8080
Restart=always
RestartSec=5
Environment=NODE_ENV=development
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target