# services/systemd/knowable-backend.service
[Unit]
Description=Knowable Backend Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/dmorris/projects/knowable/server
ExecStart=npm run dev
Restart=always
RestartSec=5
TimeoutStopSec=10
KillMode=mixed
Environment=NODE_ENV=development
Environment=PORT=3000
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target